# config/supervisor-server.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /code/config/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

[program:app-setup]
command = /bin/bash /code/config/setup_app.sh
