# config/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /code/config/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

[program:app-setup]
command = /bin/bash /code/config/setup_app.sh

[program:cron]
autostart = false
autorestart = true
command = /bin/bash -c  "declare -p | grep -Ev 'BASHOPTS|BASH_VERSINFO|EUID|PPID|SHELLOPTS|UID' > /code/cron.env && cron -f"

[program:livecheck]
autostart = false
command = /bin/bash "cd /code/app && python3 manage.py run-full-livecheck"
